# recipes-kernel/qualcomm-qca9377/files/81-wlan0-hostapd.rules
# Qualcomm QCA9377
SUBSYSTEM=="net", ACTION=="add", KERNEL=="wlan0", RUN+="/bin/systemctl restart hostapd"
